# conf/supervisor.conf
# -*- conf -*-

[program:meaning]
command=/srv/sites/meaning/conf/startup
autostart=true
autorestart=true
stderr_logfile=/srv/log/meaning/supervisor_err.log
numprocs=1
